# scheduler/crontab
# Run job.py every 1 minute
* * * * * /usr/local/bin/python /app/job.py >> /var/log/job.log 2>&1

